# data/isoftapp-daemon.service.in
[Unit]
Description=iSOFT app Service

After=network.target isoft-update-daemon.service
Wants=network.target isoft-update-daemon.service

[Service]
Type=dbus
BusName=org.isoftlinux.Isoftapp
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/isoftapp-daemon
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
